# co2monitoring.service
[Unit]
Description=AirMonitoring

[Service]
User=root
Restart=always
ExecStart=/usr/local/bin/air_monitoring.py

[Install]
WantedBy=multi-user.target